# var/lib/dpkg/info/mysql-server-8.0.list
/.
/etc
/etc/apparmor.d
/etc/apparmor.d/usr.sbin.mysqld
/etc/init
/etc/init/mysql.conf
/etc/init.d
/etc/init.d/mysql
/etc/logcheck
/etc/logcheck/ignore.d.paranoid
/etc/logcheck/ignore.d.paranoid/mysql-server-8_0
/etc/logcheck/ignore.d.server
/etc/logcheck/ignore.d.server/mysql-server-8_0
/etc/logcheck/ignore.d.workstation
/etc/logcheck/ignore.d.workstation/mysql-server-8_0
/etc/logrotate.d
/etc/logrotate.d/mysql-server
/etc/mysql
/etc/mysql/debian-start
/etc/mysql/mysql.cnf
/etc/mysql/mysql.conf.d
/etc/mysql/mysql.conf.d/mysql.cnf
/etc/mysql/mysql.conf.d/mysqld.cnf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/mysql.service
/usr
/usr/share
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_mysql-8.0.py
/usr/share/doc
/usr/share/doc/mysql-server-8.0
/usr/share/doc/mysql-server-8.0/README.Debian.gz
/usr/share/doc/mysql-server-8.0/copyright
/usr/share/doc/mysql-server-8.0/examples
/usr/share/doc/mysql-server-8.0/examples/daemon_example.ini
/usr/share/doc/mysql-server-8.0/mysqld.sym.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/mysql-server-8.0
/usr/share/mysql
/usr/share/mysql/debian_create_root_user.sql
/usr/share/mysql/docs
/usr/share/mysql/docs/INFO_BIN
/usr/share/mysql/docs/INFO_SRC
/usr/share/mysql/echo_stderr
/usr/share/mysql/install_rewriter.sql
/usr/share/mysql/mysql-systemd-start
/usr/share/mysql/mysqld_multi.server
/usr/share/mysql/uninstall_rewriter.sql
/var
/var/lib
/var/lib/mysql-upgrade
/usr/share/doc/mysql-server-8.0/NEWS.Debian.gz
/usr/share/doc/mysql-server-8.0/changelog.Debian.gz
